# aria.conf
[program:aria2server]
command=/bin/bash /usr/local/bin/aria2server.sh
autostart=true
autorestart=unexpected
stderr_logfile=/var/log/aria2server.err.log
stdout_logfile=/var/log/aria2server.out.log
